# systemd/llm-usage-indicator-tray.service
[Unit]
Description=LLM Usage Indicator Tray
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/llm-usage-indicator-tray
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
